# tools/cursor-heartbeat/com.cursor.heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd agent for Cursor Heartbeat.

  Install:
    cp com.cursor.heartbeat.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.cursor.heartbeat.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.cursor.heartbeat.plist

  Check status:
    launchctl list | grep heartbeat

  Logs:
    tail -f /tmp/heartbeat-stdout.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cursor.heartbeat</string>

    <key>ProgramArguments</key>
    <array>
        <!-- UPDATE THIS PATH -->
        <string>/Users/eli.swanson/Projects/solutions_accelerators/industry_solutions_dbxmetagen/dbxmetagen/tools/cursor-heartbeat/run_heartbeat.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>900</integer> <!-- 15 minutes -->

    <key>StandardOutPath</key>
    <string>/tmp/heartbeat-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/heartbeat-stderr.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
